# com.ga-mssr.live.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ga-mssr.live</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/karloestrada/algos/ga-mssr/.venv/bin/python</string>
        <string>/Users/karloestrada/algos/ga-mssr/scripts/run_live.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/karloestrada/algos/ga-mssr</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/karloestrada/algos/ga-mssr/logs/launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/karloestrada/algos/ga-mssr/logs/launchd_stderr.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>
</dict>
</plist>
